# services/elasticsearch/run
#!/usr/bin/env bash

cd /data
/elasticsearch/bin/elasticsearch
